# ansible/roles/hls-segmenter/templates/hls-segmenter.service.j2
[Unit]
Description=MHBN HLS segmenter for channel {{ item.id }}
After=network.target icecast2.service
Requires=icecast2.service

[Service]
User={{ hls_service_user }}
Group={{ hls_service_group }}
Restart=always
RestartSec=5
ExecStartPre=/bin/mkdir -p {{ hls_output_path }}/{{ item.output_dir }}
ExecStart=/usr/bin/ffmpeg -hide_banner -loglevel warning \
  -i http://{{ hls_icecast_host }}:{{ hls_icecast_port }}/{{ item.mount }}.ts \
  -c copy \
  -f hls \
  -hls_time {{ hls_time }} \
  -hls_list_size {{ hls_list_size }} \
  -hls_flags {{ hls_flags }} \
  -hls_segment_filename {{ hls_output_path }}/{{ item.output_dir }}/seg%%05d.ts \
  {{ hls_output_path }}/{{ item.output_dir }}/index.m3u8

[Install]
WantedBy=multi-user.target
